# packaging/systemd/airgap-downstream.service
[Unit]
Description=Air Gap Downstream Service - UDP to Kafka Receiver
Documentation=https://github.com/anders-wartoft/air-gap
Documentation=file:///usr/share/doc/airgap-downstream/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=airgap
Group=airgap

# Working directory for relative paths in config
WorkingDirectory=/etc/airgap

# Main execution
ExecStart=/usr/local/bin/airgap-downstream /etc/airgap/downstream.properties

# Environment file (optional overrides)
EnvironmentFile=-/etc/default/airgap-downstream

# Restart behavior
Restart=on-failure
RestartSec=10

# Resource limits
LimitNOFILE=65536
MemoryMax=1G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=airgap-downstream

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/airgap/downstream /var/log/airgap

[Install]
WantedBy=multi-user.target
